# docker/root/etc/services.d/20-redis/run
#!/usr/bin/execlineb -P

echo "Starting redis..."
redis-server /etc/redis.conf
