# docker/rootfs/etc/services.d/crond/run
#!/usr/bin/execlineb -P

/usr/sbin/crond -fS