# docker/root/etc/services.d/21-zealot/run
#!/usr/bin/with-contenv sh

echo "Starting zealot..."
exec /app/bundle exec puma -C /config/puma.rb