# docker/rootfs/etc/services.d/job/run
#!/usr/bin/with-contenv sh

echo "Starting Zealot Job worker"
cd /app
./bin/good_job 2>&1 | tee -a /app/log/job.log
