# docker/rootfs/etc/services.d/zealot/run
#!/usr/bin/with-contenv sh

echo "Starting Zealot"
cd /app
bundle exec puma -C config/puma.rb  | tee /app/log/zealot.log