# docker/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv sh

echo "Starting Caddy"
/usr/sbin/caddy run --config /etc/Caddyfile > /dev/null 2>&1
