# docker/rootfs/etc/services.d/sidekiq/run
#!/usr/bin/with-contenv sh

echo "Starting Zealot Worker"
/app/bin/sidekiq -C /app/config/sidekiq.yml  | tee /app/log/sidekiq.log
