# docker/rootfs/etc/services.d/sidekiq/run
#!/usr/bin/with-contenv sh

echo "Starting Zealot Worker"
cd /app
./bin/sidekiq -C config/sidekiq.yml  | tee /app/log/sidekiq.log
